# files/package/base-files/files/etc/hotplug.d/iface/99-netfix
#!/bin/sh
case "$ACTION" in
ifup)
/sbin/video.sh
/etc/init.d/redsocks2 enabled && /etc/init.d/redsocks2 restart
;;
ifdown)
;;
esac
