# files/package/base-files/files/etc/rc.local
/etc/init.d/redsocks enable
chmod 600 /etc/stunnel/stunnel.pem
